# docker/php/crontab
*/15 * * * * /usr/local/bin/php /var/www/html/bot.php >> /bot-cron.log 2>&1
